# _packaging/nfpm/guard.timer
[Unit]
Description=immune Guard agent timer

[Timer]
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target
